# docs/files/storage-thin.conf
description "thin http daemon"
author      "Peter A Leonov <gojpeg@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]


env LANG=en_US.UTF-8
exec thin start -a 127.0.0.1 -p 4000 -u www -g www -R "/www/storage/thin.ru" -e production
# exec thin start -a 127.0.0.1 -p 4000 -u www -g www -R "/www/storage/thin.ru" -P "/www/storage/logs/thin.pid" -l "/www/storage/logs/thin.log" --tag storage


respawn
respawn limit 60000 1
