# storage/upstart/storage-nginx.conf
description "storage's nginx"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/nginx/sbin/nginx -p "/www/storage/" -c "nginx.production.conf"

respawn
respawn limit 60000 1
