# storage/upstart/storage-thin.conf
description "storage's thin"

start on runlevel [2345]
stop on runlevel [!2345]

env LANG=en_US.UTF-8

exec thin start -a 127.0.0.1 -p 18183 -u www -g www -R "/www/inshaker/storage/thin.ru" -l "/www/inshaker/storage/logs/thin.log"

respawn
respawn limit 60000 1
